# platform/scripts/crontab
*/5 * * * *   root    /srv/beialand/platform/scripts/git_update.sh
